# target/linux/pxa/zipitz2/base-files/etc/init.d/z2-sound
#!/bin/sh /etc/rc.common
# Set Zipit Z2 Sound

START=96

start() {
	alsactl restore -f /etc/asound.state
}

